# copy_to_systemd_user/voice.service
[Unit]
Description=BG/JA voice loop
After=network-online.target sound.target

[Service]
Working Directory=%h/app
Environment="PATH=%h/venvs/ai/bin:/usr/bin"

#オフライン動作の環境変数(ここに書けば常に有効)
Environment=HF_HUB_OFFLINE=1
Environment=TRANSFORMERS_OFFLINE=1
ExecStart=%h/venvs/ai/bin/python run_voice.py

Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target